# dev-ruby/dyndyndong/files/initd
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	need net
	use logger dns
	after bootmisc
}

start() {
	ebegin "Starting DynDynDong"
	start-stop-daemon --start --quiet -c dyndyndong:dyndyndong --exec /usr/bin/dyndyndong \
		--pidfile /var/run/dyndyndong.pid --make-pidfile --background --
	eend $?
}

stop() {
	ebegin "Stopping DynDynDong"
	start-stop-daemon --stop --quiet --exec /usr/bin/dyndyndong \
		--pidfile /var/run/dyndyndong.pid
	eend $?
}
